# templates/default/unicorn.erb
#!/bin/bash
# /etc/init.d/unicorn
# To install: update-rc.d unicorn defaults
# Snarfed then tweaked from https://gist.github.com/2814835

### BEGIN INIT INFO
# Provides:   unicorn
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Should-Start:      $network
# Should-Stop:       $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Unicorn server
# Description:       Starts the unicorn server
### END INIT INFO

set -u
set -e
# Example init script, this can be used with nginx, too,
# since nginx and unicorn accept the same signals

TIMEOUT=${TIMEOUT-60}

CMD="unicorn-<%= @env %>"

cd <%= @root %> || exit 1

sig() {
  test -s "<%= @pid %>" && kill -$1 `cat <%= @pid %>`
}

old_pid="<%= @pid %>.oldbin"

oldsig () {
  test -s $old_pid && kill -$1 `cat $old_pid`
}

start() {
  su -l -c '<%= @command %>' <%= @user %>
}

case $1 in
start)
  sig 0 && echo >&2 "Already running" && exit 0
  start
  ;;
stop)
  sig QUIT && exit 0
  echo >&2 "Not running"
  ;;
force-stop)
  sig TERM && exit 0
  echo >&2 "Not running"
  ;;
restart|reload)
  sig HUP && echo reloaded OK && exit 0
  echo >&2 "Couldn't reload, starting unicorn-<%= @env %> instead"
  start
  ;;
status)
  test -s "<%= @pid %>" && test -d "/proc/`cat <%= @pid %>`" && echo running && exit 0
  echo stopped && exit 3
  ;;
upgrade)
  if sig USR2 && sleep 2 && sig 0 && oldsig QUIT
  then
    n=$TIMEOUT
    while test -s $old_pid && test $n -ge 0
    do
      printf '.' && sleep 1 && n=$(( $n - 1 ))
    done
    echo

    if test $n -lt 0 && test -s $old_pid
    then
      echo >&2 "$old_pid still exists after $TIMEOUT seconds"
      exit 1
    fi
    exit 0
  fi
  echo >&2 "Couldn't upgrade, starting unicorn-<%= @env %> instead"
  start
  ;;
rotate)
  sig USR1 && echo rotated logs OK && exit 0
  echo >&2 "Couldn't rotate logs" && exit 1
  ;;
*)
  echo >&2 "Usage: $0 <start|stop|restart|upgrade|rotate|force-stop>"
  exit 1
  ;;
esac
